# QuizMaster/launchd/com.easypv.quiz-helper-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.easypv.quiz-helper-watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/thor/JetBrains/WebstormProjects/easyPV/databases/entwicklung/project_tickets_v01/QuizMaster/scripts/helper_watchdog_codex.py</string>
    <string>--once</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/thor/JetBrains/WebstormProjects/easyPV/databases/entwicklung/project_tickets_v01</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/thor</string>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Applications/Codex.app/Contents/Resources</string>
    <key>CODEX_HOME</key>
    <string>/Users/thor/.codex</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>600</integer>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Users/thor/JetBrains/WebstormProjects/easyPV/databases/entwicklung/project_tickets_v01/QuizMaster/output/helper_watchdog_launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/thor/JetBrains/WebstormProjects/easyPV/databases/entwicklung/project_tickets_v01/QuizMaster/output/helper_watchdog_launchd.err.log</string>
</dict>
</plist>
